# src/assembly/debian11/control/prerm
#!/usr/bin/env bash

DEBIAN_NAME="dc-con"

/bin/systemctl stop $DEBIAN_NAME.service
/bin/systemctl disable $DEBIAN_NAME.service

exit 0