# src/assembly/debian9/control/prerm
#!/bin/bash

/bin/systemctl stop dc-con.service
/bin/systemctl disable dc-con.service

exit 0