# src/assembly/debian9/control/prerm
#!/bin/bash

/bin/systemctl stop smb-con.service
/bin/systemctl disable smb-con.service

exit 0